# infra/server-green.service
[Unit]
Description=Gleam Web Server (artisan green)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/artisan/green/server
Environment=PORT=8001
ExecStart=/usr/local/bin/gleam run
Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

