# infra/server.service
[Unit]
Description=Gleam Web Server (artisan/server)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/artisan/server
ExecStart=/usr/local/bin/gleam run
Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

