# infra/server-green.service
[Unit]
Description=Gleam Web Server (artisan green)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/artisan/green/server
Environment=PORT=8001
Environment=PATH={{ERLANG_BIN_DIR}}:/usr/bin:/bin
ExecStart={{GLEAM_BIN}} run
Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

